# systemd/linkedin-cron-scheduler.service
[Unit]
Description=LinkedIn Cron Scheduler Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/opt/linkedin-cron
EnvironmentFile=/etc/linkedin-cron.env
ExecStart=/usr/local/bin/linkedin-cron-scheduler
StandardOutput=journal
StandardError=journal
