# systemd/stroopwafel.service
[Unit]
Description=Stroopwafel: Social Media Manager HTTP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/stroopwafel
EnvironmentFile=/etc/stroopwafel.env
ExecStart=/usr/local/bin/stroopwafel-server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
